# user/gnome-keyring-init.path
[Unit]
Description=GNOME Keyring Initializer
Documentation=man:gnome-keyring-daemon(1)

[Path]
PathExists=%t/keyring/control

[Install]
WantedBy=default.target
