# user/x11-session-pre.target
[Unit]
Description=X11 Window System Pre-Session

BindsTo=graphical-session-pre.target
Before=graphical-session-pre.target
